# deploy/systemd/auditoria-m2m.service
[Unit]
Description=Auditoria M2M Streamlit
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/auditmem
Environment="PATH=/opt/auditmem/.venv/bin"
ExecStart=/opt/auditmem/.venv/bin/streamlit run app.py --server.address 127.0.0.1 --server.port 8502
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
